# sd-webui/sd-webui/sd-webui.conf
[program:sd-webui]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/sd-webui.sh
user=root
directory=/root/stable-diffusion-webui
