# wan2.1t2v/wan2.1t2v/wan2.1t2v.conf
[program:wan2.1t2v]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/wan2.1t2v.sh
user=root
directory=/root
stopasgroup=true
killasgroup=true
