# fluxgym/fluxgym/fluxgym.conf
[program:fluxgym]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/fluxgym.sh
user=root
directory=/app/fluxgym
stopasgroup=true
killasgroup=true
