# wan2.1t2v/comfyui/comfyui.conf
[program:comfyui]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/comfyui.sh
user=root
directory=/root
stopasgroup=true
killasgroup=true
