# HunyuanVideo/hunyuanvideo/hunyuanvideo.conf
[program:hunyuanvideo]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/hunyuanvideo.sh
user=root
directory=/root
stopasgroup=true
killasgroup=true
