# HunyuanVideo/sshd/sshd.conf
[program:sshd]
autorestart=True
autostart=True
environment=
command=/bin/bash /opt/supervisor/sshd.sh
user=root
directory=/root
